# build/fogproject-install.service
[Unit]
Description=FOG Project Install Script Runner
After=network.target

[Service]
Type=idle
ExecStart=/bin/true
ExecStartPre=/install-fogproject.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
